# signal_review/launchd/com.nexus.signal-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nexus.signal-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/raymonddavis/nexus-ai/scripts/start_signal_review.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/raymonddavis/nexus-ai</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>15</integer>

  <key>StandardOutPath</key>
  <string>/Users/raymonddavis/nexus-ai/openclaw/logs/signal-review.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/raymonddavis/nexus-ai/openclaw/logs/signal-review.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/raymonddavis</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
